# services/cah.timer
[Unit]
Description = Fetch Calvin & Hobbes comic

[Install]
WantedBy = timers.target

[Timer]
OnCalendar=daily
OnActiveSec=180
Persistent=true
